# config/aphlict.sv.conf
[program:aphlict]
command=/srv/phabricator/phabricator/bin/aphlict debug --config /srv/phabricator/phabricator/conf/aphlict/aphlict.default.json
user=phabricator-daemon